# sync-settings/files/interface-watcher.init
#!/bin/sh /etc/rc.common

# 16 is right after startup (15)
START=19
STOP=19

USE_PROCD=1

start_service() {
	procd_open_instance interface-watcher
	procd_set_param command /usr/bin/interface_watcher.sh
	procd_close_instance
}
