# sync-settings/files/netlink-detect.init
#!/bin/sh /etc/rc.common

# 17 is right after setup (16)
START=18
STOP=18

start() {
    if [ -f /etc/Eos-release ] ; then
        netlink-detect
    fi
}

stop() {
    # Not implemented
    true
}
