# sync-settings/files/setup.init
#!/bin/sh /etc/rc.common

# 15 is right after sync-settings (14)
START=16
STOP=16

boot() {
    if [ -d /etc/config/setup.d/ ] ; then
       	run-parts /etc/config/setup.d/
	fi

	if [ -f /etc/config/zzz-reload-system ] ; then
		/etc/config/zzz-reload-system
	fi
}
