# sync-settings/files/static-routes.init
#!/bin/sh /etc/rc.common

# Needs to be started after network!
START=99
STOP=99

boot() {
	ubus call network restart
}
