# sync-settings/files/static-routes.init
#!/bin/sh /etc/rc.common

# Needs to be started after network!
START=99
STOP=99

start() {
	true 
	echo "NO restart of network" >>/tmp/test.log
	#ubus call network restart
}
stop() {
	# Not implemented
	true
}