# sync-settings/files/test-detect.init
#!/bin/sh /etc/rc.common

# 17 is right after setup (16)
START=18
STOP=18

start() {
   ls -ltr /etc/init.d > /etc/config/my.txt
}

stop() {
    # Not implemented
    true
}
