# license-scripts/files/fetch-licenses.init
#!/bin/sh /etc/rc.common

# this can happen late
START=99
STOP=99

boot() {
    # if we dont' have a license file on boot
    # try fetching one
    if [ ! -f /etc/config/licenses.json ] ; then
        /usr/bin/fetch-licenses.sh
    fi
}
