# build/sync-settings/files/uid.init
#!/bin/sh /etc/rc.common

# 11 is right after boot (10)
START=11
STOP=11

boot() {
	if [ ! -f /etc/config/uid ] ; then
	    uuidgen -r > /etc/config/uid
	fi
}
